# com.user.icloud_smb_backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.icloud_smb_backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/{{USERNAME}}/iCloud Backup/icloud_smb_backup.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/{{USERNAME}}/Library/Mobile Documents/com~apple~CloudDocs/</string>
    </array>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/com.user.icloud_smb_backup.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.user.icloud_smb_backup.stderr.log</string>
</dict>
</plist>
